# test/supervisor-docker.conf
# Remote CI Supervisor 配置 - Docker 版本
# 所有服务在同一个容器中运行

[program:redis]
command=redis-server --daemonize no --appendonly yes --dir /app/data
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/remote-ci/redis.log
stderr_logfile=/var/log/remote-ci/redis.log
priority=10

[program:api]
command=python -m server.app
directory=/app
environment=PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/remote-ci/api.log
stderr_logfile=/var/log/remote-ci/api.log
priority=20

[program:worker]
command=celery -A server.celery_app worker --loglevel=info --concurrency=2
directory=/app
environment=PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/remote-ci/worker.log
stderr_logfile=/var/log/remote-ci/worker.log
stopwaitsecs=60
stopasgroup=true
killasgroup=true
priority=30

[program:flower]
command=celery -A server.celery_app flower --port=5555
directory=/app
environment=PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
autostart=false
autorestart=true
startretries=3
stdout_logfile=/var/log/remote-ci/flower.log
stderr_logfile=/var/log/remote-ci/flower.log
priority=40

[group:remote-ci]
programs=redis,api,worker
priority=999
